# debian/postinst
#!/bin/sh

useradd --system simeis || true

systemctl daemon-reload
systemctl enable simeis
systemctl start simeis