# debian/prerm
#!/bin/sh

systemctl stop simeis
systemctl disable simeis